# keylogger.service
[Unit]
Description=keylogger

[Service]
ExecStart=/home/caente/bin/keylogger -f /home/caente/keys.log

[Install]
WantedBy=multi-user.target
